# .config/systemd/user/btpoweron.service
[Unit]
Description=Power on bluetooth.

[Service]
Type=simple
ExecStart=/usr/bin/bluetoothctl power on

[Install]
WantedBy=default.target
